# android_emulator/set-black-bg.conf
[program:set-black-bg]
command=/bin/sh -c "sleep 10 && DISPLAY=:0 xsetroot -solid black"
autorestart=false
startretries=1
priority=999
stdout_logfile=/dev/null
stderr_logfile=/dev/null
